# shell/asp.ini
[unix_http_server]
file=/app/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
childlogdir=/var/log/supervisor
logfile=/app/log/supervisord.log
pidfile=/app/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///app/supervisor.sock

[program:crond]
command=crond -f
autostart=true
autorestart=true
stdout_logfile=/app/log/ASP.log