# app/asp.conf
[unix_http_server]
file=/app/supervisor.sock
chmod=0700

[supervisord]
nodaemon=true
childlogdir=%(ENV_CONFIG_DIR)s/logs/supervisor
logfile=%(ENV_CONFIG_DIR)s/logs/supervisord.log
pidfile=/app/supervisord.pid
user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///app/supervisor.sock

[program:crond]
command=crond -f
autostart=true
autorestart=true
stdout_logfile=%(ENV_CONFIG_DIR)s/logs/ASP.log